# roles/caddy/templates/modify_nginx.timer.j2
[Unit]
Description=Make sure NGiNX configuration is modified after boot
After=modify_nginx.path
Wants=modify_nginx.path

[Timer]
OnBootSec=3min
